# .local/sv/rivertile/run
#!/bin/sh
exec rivertile -view-padding 2 -outer-padding 2 -main-ratio 0.50
