# .local/sv/pipewire/run
#!/bin/sh
exec pipewire -c /home/robin/.config/pipewire/pipewire.conf
